# l.service
[Unit]
Description=Lighter Trading Bot - Симуляция торговли на Lighter бирже
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/ser_lighter
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=-/root/ser_lighter/.env

# Команда запуска
ExecStart=/usr/bin/python3 /root/ser_lighter/lighter_trading_bot.py

# Перезапуск при сбое
Restart=always
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Логирование
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lighter-trading-bot

# Безопасность
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/root/ser_lighter

# Ограничения ресурсов (опционально)
# MemoryMax=2G
# CPUQuota=100%

[Install]
WantedBy=multi-user.target

